# gentoo/app-emulation/horcrux/files/horcrux-api.initd
#!/sbin/openrc-run
# Copyright 2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="Horcrux API Server"
description="Virtualization management API server"

command="/usr/bin/horcrux-api"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need net
	use logger
}
